# debian/renasar-syslog.upstart
description "run the renasar syslog service"

env config=/etc/default/renasar-syslog

start on started rabbitmq-server
stop on runlevel [!2345]
respawn

pre-start script
[ -e "$config" ] || { echo "${config} missing, not starting"; stop; exit 0; }
end script

script
    echo Starting renasar syslog service
    chdir /var/renasar/renasar-syslog
    exec /usr/bin/nodejs /var/renasar/renasar-syslog/index.js
end script
